# root/etc/services.d/gollum-svc/run
#!/usr/bin/with-contenv bash
cd /wiki
gollum --port 8080